# scripts/systemd/sensor-heartbeat.service
[Unit]
Description=Sensor System Heartbeat Service (C++ — SERVER_HEARTBEAT to boards)
After=network.target sensor-backend.service

[Service]
Type=simple
WorkingDirectory=%h/sensor_system
ExecStart=/bin/bash -c 'if [ -x ./build/FSW/heartbeat_service ]; then exec ./build/FSW/heartbeat_service --config config/config.toml; elif [ -x ./FSW/build/heartbeat_service ]; then exec ./FSW/build/heartbeat_service --config config/config.toml; else exec python3 scripts/services/heartbeat_service.py --config config/config.toml; fi'
Restart=always
RestartSec=5
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=default.target
