# deploy/systemd/sensor-controller.service
[Unit]
Description=Sensor System Controller Service (C++ — Elodin→RobustDDP→PWM)
After=network.target sensor-elodin.service sensor-relay.service

[Service]
Type=simple
WorkingDirectory=%h/sensor_system
ExecStart=/bin/bash -c 'if [ -x ./build/bin/controller_service ]; then exec ./build/bin/controller_service --config config/config.toml --elodin-host 127.0.0.1; else echo "controller_service not found — run: cd build && cmake .. && make controller_service"; exit 1; fi'
Restart=always
RestartSec=5
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=default.target
