# daq-server/scripts/systemd/sensor-daq.service
[Unit]
Description=Sensor System DAQ Bridge (UDP → Elodin DB)
After=network.target sensor-elodin.service

[Service]
Type=simple
WorkingDirectory=%h/sensor_system
ExecStart=/bin/bash -c 'if [ -x ./build/FSW/daq_bridge ]; then exec ./build/FSW/daq_bridge config/config.toml; elif [ -x ./FSW/build/daq_bridge ]; then exec ./FSW/build/daq_bridge config/config.toml; else echo "daq_bridge not found — run: cd FSW/build && cmake ../.. && make daq_bridge"; exit 1; fi'
Restart=always
RestartSec=5
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=default.target
